# examples/clu.discord_inbound.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent template for the clu Discord inbound poller.

  Install:
    1. Copy to ~/Library/LaunchAgents/com.clu.discord_inbound.plist
    2. Edit ProgramArguments: replace /usr/bin/python3 with the python that
       has end_of_line installed (e.g. `which python3` from the venv where
       you ran `pip install -e .`). LaunchAgents do NOT inherit your shell
       PATH — full paths are mandatory.
    3. Edit the --project argument to point at your project root (the
       directory containing .orchestrator.json).
    4. launchctl load ~/Library/LaunchAgents/com.clu.discord_inbound.plist
       (or `launchctl bootstrap gui/$UID …` on macOS 10.11+)
    5. Tail /tmp/clu-discord-inbound.{out,err} to confirm the loop is running.

  Prerequisites:
    - A Discord bot token and your Discord user ID configured in
      .orchestrator.json under notify.channels (kind: "discord").
    - bot_user_id also set in the channel spec (the bot's own user ID,
      used to filter the bot's outbound messages from the inbound stream).

  KeepAlive=true restarts the poller if it crashes; ThrottleInterval
  prevents a crash-loop from saturating launchd.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clu.discord_inbound</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>end_of_line.notify_discord_inbound</string>
        <string>--project</string>
        <string>/path/to/your/project</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/clu-discord-inbound.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/clu-discord-inbound.err</string>
</dict>
</plist>
